# server/systemd/device-policy-manager.service.in
[Unit]
Description=Device Policy Manager

[Service]
Type=simple
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/device-policy-manager
Restart=on-failure
ExecReload=/bin/kill -HUP $MAINPID
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE

[Install]
WantedBy=multi-user.target
